# valkyrie/systemd/valkyrie-worker.service
[Unit]
Description=Valkyrie Worker (Docker Compose)
After=docker.service
Requires=docker.service

[Service]
# Load environment variables from user's valkyrie .env so host-port and other overrides apply.
EnvironmentFile=-%h/valkyrie/.env
WorkingDirectory=%h/valkyrie
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/docker compose -f docker-compose.yml up -d --build worker
ExecStop=/usr/bin/docker compose -f docker-compose.yml stop worker || /usr/bin/docker compose -f docker-compose.yml down
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
